# pozz-monitor.service
[Unit]
Description=POZZ Application Monitor
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/path/to/POZZ
ExecStart=/path/to/POZZ/monitor.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="CHECK_INTERVAL=300"
Environment="MAX_RESTART_ATTEMPTS=5"
Environment="PORT=8501"

[Install]
WantedBy=multi-user.target


